# .config/systemd/user/cli-proxy-api.service
[Unit]
Description=CLI Proxy API Service
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/cli-proxy-api --config %h/.config/cli-proxy-api.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
