# .config/systemd/user/psd-resync.timer
[Unit]
Description=Timer for profile-sync-daemon - 25min
BindsTo=psd.service

[Timer]
OnUnitActiveSec=25m
