# .config/systemd/user/bitwarden-keyring.service.d/noctalia.conf
[Service]
Environment=BITWARDEN_KEYRING_NOCTALIA=1
